# docker/conf/supervisord.ini
 
[supervisord]
nodaemon=true

[program:syslog]
command=rsyslogd -n -f /etc/rsyslog.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=QUIT

[program:sd-test-srv]
command=sd-test-srv --listen ":53"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=syslog:running
autostart=false
stopsignal=QUIT

[program:boulder-remoteva]
command=boulder-remoteva --config /opt/boulder/config/va-remote-a.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running syslog:running
autostart=false
stopsignal=QUIT


[program:boulder-sa]
command=boulder-sa --config /opt/boulder/config/sa.json --addr sa.$INTERNAL_DOMAIN:9095
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running syslog:running
autostart=false
stopsignal=QUIT

[program:ct-test-srv]
command=ct-test-srv --config /opt/boulder/config/ct-test-srv/ct-test-srv.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=syslog:running
autostart=false
stopsignal=QUIT

[program:boulder-publisher]
command=boulder-publisher --config /opt/boulder/config/publisher.json --addr "publisher.$INTERNAL_DOMAIN:9091" --debug-addr', ":8009"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running syslog:running
autostart=false
stopsignal=QUIT

[program:mail-test-srv]
command=mail-test-srv --closeFirst 5 --cert /opt/boulder/config/secrets/ca/mail-test-srv/cert.pem --key /opt/boulder/config/secrets/ca/mail-test-srv/key.pem
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=syslog:running
autostart=false
stopsignal=QUIT

[program:ocsp-responder]
command=ocsp-responder --config /opt/boulder/config/ocsp-responder.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=syslog:running
autostart=false
stopsignal=QUIT

[program:boulder-va]
command=boulder-va --config /opt/boulder/config/va.json --addr va.$INTERNAL_DOMAIN:9092 --debug-addr ":8004"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running boulder-remoteva:running  syslog:running
autostart=false
stopsignal=QUIT

[program:boulder-ca]
command=boulder-ca  --config /opt/boulder/config/ca-a.json --ca-addr ca.$INTERNAL_DOMAIN:9093 --ocsp-addr ca.$INTERNAL_DOMAIN:9096 --debug-addr ":8001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running boulder-sa:running syslog:running
autostart=false
stopsignal=QUIT

[program:akamai-test-srv]
command=akamai-test-srv --listen localhost:6789 --secret "$AKAMAI_SECRET" 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=syslog:running
autostart=false
stopsignal=QUIT

[program:akamai-purger]
command=akamai-purger --config /opt/boulder/config/akamai-purger.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running akamai-test-srv:running syslog:running
autostart=false
stopsignal=QUIT

[program:ocsp-updater]
command=ocsp-updater --config /opt/boulder/config/ocsp-updater.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running boulder-ca:running syslog:running
autostart=false
stopsignal=QUIT

[program:boulder-ra]
command=boulder-ra --config /opt/boulder/config/ra.json --addr ra.$NTERNAL_DOMAIN:9094 --debug-addr ":8002"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running boulder-sa:running boulder-ca-a:running boulder-va:running akamai-purger:running boulder-publisher:running syslog:running
autostart=false
stopsignal=QUIT

[program:bad-key-revoker]
command=bad-key-revoker --config /opt/boulder/config/bad-key-revoker.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running boulder-ra:running mail-test-srv:running syslog:running
autostart=false
stopsignal=QUIT

[program:nonce-service]
command=nonce-service --config /opt/boulder/config/nonce.json --addr nonce.$INTERNAL_DOMAIN:9101 --debug-addr :8111 --prefix taro
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running syslog:running
autostart=false
stopsignal=QUIT

[program:boulder-wfe]
command=boulder-wfe --config /opt/boulder/config/wfe.json
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
dependent_startup=true
dependent_startup_wait_for=sd-test-srv:running boulder-ra:running boulder-sa:running nonce-service:running syslog:running
autostart=false
stopsignal=QUIT
